# meta-ebaz4205/recipes-core/addr-led/files/addr-led.service
[Unit]
Description=Adressable led strip controller

[Service]
Type=simple
Restart=always
ExecStart=led-strip -n 120

[Install]
WantedBy=multi-user.target
